# .platform/files/supervisor.ini
# The following given example is used with Horizon, but you can
# make a copy of it or add your own and configure it.

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/app/current/artisan queue:work sqs --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/var/app/current/storage/worker.log
stopwaitsecs=3600